# nowchat/supervisord.conf
[program:nowchat]
command = node index.js
directory = /home/dotcloud/current
